# backups/backups.cron
# Cron file for running backup operations

# run backups every 5 hours
0 */5 * * * /app/backup-db.sh
